# backend/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:fastapi]
command=uvicorn app.main:app --host 0.0.0.0 --port 8000 --log-level info
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/fastapi.err.log
stdout_logfile=/var/log/supervisor/fastapi.out.log
environment=PYTHONUNBUFFERED=1

[program:pinchtab]
command=pinchtab daemon start
directory=/app
autostart=true
autorestart=true
startsecs=5
stderr_logfile=/var/log/supervisor/pinchtab.err.log
stdout_logfile=/var/log/supervisor/pinchtab.out.log
# PinchTab is optional — if it fails, FastAPI still runs
exitcodes=0,1,2
